# supervisord-terminal.conf
[program:terminal]
command=/start-terminal.sh
numprocs=1
autostart=true
autorestart=true